# deploy/roles/vnc/templates/x11vnc.conf
start on login-session-start 

script
  x11vnc -display :0 -auth /var/run/lightdm/root/:0 -forever -bg -o /var/log/x11vnc.log -rfbauth /etc/x11vnc.passwd -rfbport 5900
end script
